# aio/supervisord.conf
[program:webserver]
command=gunicorn -b '127.0.0.1:{{ settings.SERVER_PORT }}' --workers 3 wsgi:application --error-logfile /home/futuregreen/spontaneousurbanplants.org/logs/gunicorn/error.log
